# dot_config/systemd/user/mbsync.service.tmpl
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart={{ env "SHELL" }} -c check-mail
